# fetch_feeds.cron
0 * * * * /ug/ug2k6/cse/gaganpreet/public_html/blogroll/fetch_feeds.sh
